# rootfs/etc/s6/skydns/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

STARTCMD="/usr/bin/skydns -addr=0.0.0.0:53 -systemd=false"


exec ${STARTCMD}
